# rootfs/etc/rc.local
#!/bin/bash

# mmc led
echo mmc1 > /sys/devices/soc0/gpio-leds/leds/wifi-led/trigger

# set brightness, not need really while use screen blanking
# echo 2 > /sys/devices/soc0/backlight/backlight/backlight/brightness

# boost CPU
echo performance > /sys/devices/system/cpu/cpufreq/policy0/scaling_governor

# tuning cache
echo 40 > /proc/sys/vm/dirty_background_ratio
echo 10000 > /proc/sys/vm/dirty_writeback_centisecs
echo 9999999 > /proc/sys/vm/dirty_background_bytes
echo 90 > /proc/sys/vm/dirty_background_ratio
echo 30000 > /proc/sys/vm/dirty_expire_centisecs
echo 300000 > /proc/sys/vm/dirty_writeback_centisecs
echo 100 >  /proc/sys/vm/swappiness

# DPMS feature works for screen powersaving
TERM=linux setterm --blank 1 --msg on  >/dev/tty1 </dev/tty1

# you may set a splash directly to framebuffer
#cat /etc/default/toshibaAC100.fbimg  > /dev/fb0
exit 0
